# PDP-8/pidp8i/etc/usb-mount@.service.in
[Unit]
Description=Mount USB Drive on %i
[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=@ABSPREFIX@/bin/usb-mount add %i
ExecStop =@ABSPREFIX@/bin/usb-mount remove %i
